# configs/systemd/tallow-app.service
# Systemd service file for Tallow App
# Place this file in: /etc/systemd/system/tallow-app.service
# Enable with: sudo systemctl enable tallow-app
# Start with: sudo systemctl start tallow-app

[Unit]
Description=Tallow Web Application - Secure P2P File Sharing
Documentation=https://github.com/yourusername/tallow
After=network.target network-online.target
Wants=network-online.target tallow-signaling.service

[Service]
Type=simple
User=tallow
Group=tallow
WorkingDirectory=/var/www/tallow

# Environment
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0
EnvironmentFile=/var/www/tallow/.env.local

# Main process
ExecStart=/usr/bin/npm start

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Performance tuning
LimitNOFILE=65535
Nice=-10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/tallow
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tallow-app

[Install]
WantedBy=multi-user.target
